# submission/week1/07.service
[Unit]
Description= Wolex Service

[Service]
Type=oneshot
ExecStart=/bin/echo "Hello World"