# submission/week1/07.timer
[Unit]
Description=Wolex Timer

[Timer]
Unit=07.service
OnCalendar=minutely

[Install]
WantedBy=timers.target